# deploy/systemd/hostel-next.service
[Unit]
Description=Hostel Management Next.js
After=network.target

[Service]
WorkingDirectory=/opt/hostel-management/frontend
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
