# deploy/systemd/hostel-worker.service
[Unit]
Description=Hostel Management Notification Worker
After=network.target

[Service]
WorkingDirectory=/opt/hostel-management
Environment="PATH=/opt/hostel-management/.venv/bin"
EnvironmentFile=/opt/hostel-management/.env
ExecStart=/opt/hostel-management/.venv/bin/python -m worker.run_notifications_loop
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
